# history/2015-11-15/worker/rc.local
sed -i "s/numprocs=.*/numprocs=$(nproc)/" /etc/supervisor/conf.d/worker.conf

supervisorctl reload || :
