# history/2015-11-15/cws/sup_cws.conf
[program:cmsContestWebServer-jau]
command=/usr/local/bin/cmsContestWebServer -c 1 0
user=ubuntu

[program:cmsContestWebServer-vyr]
command=/usr/local/bin/cmsContestWebServer -c 2 1
user=ubuntu
