# history/2015-11-15/centriukas/sup_cms.conf
[program:cmsLogService]
command=/usr/local/bin/cmsLogService
user=ubuntu

[program:cmsResourceService]
command=/usr/local/bin/cmsResourceService -a 1
user=ubuntu
stopasgroup=true
